# setup/home_intent.conf
[program:home_intent]
command=/usr/bin/python3 /usr/src/app/home_intent
stopasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
